# deploy/supervisor/asas-horizon.conf
; =====================================================================
; Supervisor program — Laravel Horizon (ASAS Investor Portal)
; Place at: /etc/supervisor/conf.d/asas-horizon.conf  (Linux production)
; Then:  sudo supervisorctl reread && sudo supervisorctl update
;        sudo supervisorctl start asas-horizon
; Horizon requires the pcntl + posix PHP extensions (Linux only).
; =====================================================================
[program:asas-horizon]
process_name=%(program_name)s
command=php /var/www/asas-portal/artisan horizon
directory=/var/www/asas-portal
autostart=true
autorestart=true
; Give Horizon time to finish in-flight jobs on a graceful restart.
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/asas-portal/storage/logs/horizon.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

; -----------------------------------------------------------------
; Deploys should gracefully restart Horizon (drains current jobs):
;     php artisan horizon:terminate
; Supervisor then auto-restarts it with the new code.
; -----------------------------------------------------------------
